# enstore-transfer/gratia-probe-enstore-transfer.cron
# To send the output to a specific email add MAILTO=user@host
10,40 * * * * root /usr/share/gratia/common/cron_check  /etc/gratia/enstore-transfer/ProbeConfig && /usr/share/gratia/enstore-transfer/enstore-transfer -c -s 1800
