# enstore-transfer/gratia-probe-enstore-transfer.cron
10 5,17 * * * root /usr/share/gratia/common/cron_check  /etc/gratia/enstore-transfer/ProbeConfig && /usr/share/gratia/enstore-transfer/enstore-transfer -c -s 1800
