# glexec/gratia-probe-glexec.cron
0 * * * * root /usr/share/gratia/common/cron_check  /etc/gratia/glexec/ProbeConfig  && perl -e 'sleep rand 3600' && /usr/share/gratia/glexec/glexec_meter.cron.sh
