# glexec/gratia-probe-glexec.cron
0 * * * * root perl -e 'sleep rand 3600' && /usr/share/gratia/glexec/glexec_meter.cron.sh
