# enstore-storage/gratia-probe-enstore-storage.cron
0,30 * * * * root /usr/share/gratia/common/cron_check  /etc/gratia/enstore-storage/ProbeConfig && /usr/share/gratia/enstore-storage/enstore-storage -s 1800
