# enstore-storage/gratia-probe-enstore-storage.cron
# Since enstore space data is currently updated once a month (23:30 of the last day) and may increase to weekly, there is no reason to run this probe more than once a day
# To send the output to a specific email add MAILTO=user@host
40 0 * * * root /usr/share/gratia/common/cron_check  /etc/gratia/enstore-storage/ProbeConfig && /usr/share/gratia/enstore-storage/enstore-storage -c -s 1800
