# psacct/gratia-probe-psacct.cron
0 1 * * * root /usr/share/gratia/common/cron_check  /etc/gratia/psacct/ProbeConfig  && perl -e 'sleep rand 9600' && /usr/share/gratia/psacct/psacct_probe.cron.sh
