# gridftp-transfer/gratia-probe-gridftp-transfer.cron
0,30 * * * * root /usr/share/gratia/gridftp-transfer/GridftpTransferProbeDriver -s 1800
